# resources/celery.cron
*/5 * * * * /Application/bermuda/monitor_celery.sh >>/Application/bermuda/logs/monitor_celery.log 2>&1